# supervisord.conf
[supervisord]
;logfile = /logs/supervisord.log
logfile_maxbytes = 0
;logfile_backups=10
logfile=/dev/stdout
loglevel = debug
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 65536
minprocs = 200
umask = 022
user = root
identifier = supervisor
directory = /app
childlogdir = /logs

[program:php-fpm]
user=root
command=/usr/local/php/sbin/php-fpm
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
;stdout_logfile=/logs/stdout.log
;stdout_logfile_maxbytes=10MB
;stderr_logfile=/logs/stderr.log
;stderr_logfile_maxbytes=10MB
redirect_stderr=true
redirect_stout=true

[program:nginx]
user = root
command=/usr/local/nginx/sbin/nginx -g "daemon off;"
;stdout_logfile=/logs/stdout.log
;stdout_logfile_maxbytes=10MB
;stderr_logfile=/logs/stderr.log
;stderr_logfile_maxbytes=10MB
redirect_stderr=true
redirect_stout=true

